# res/systemd/als-led-backlight.service
[Unit]
Description=als-led-backlight Service

[Service]
ExecStart=/opt/als-led-backlight

[Install]
WantedBy=multi-user.target